# task05/98-task05.rules
# do not edit this file, it will be overwritten on update

#ACTION=="remove", GOTO="task05_remove"
SUBSYSTEM!="usb*", GOTO="task05_end"


ACTION=="add", ATTRS{idVendor}=="04ca", ATTRS{idProduct}=="0061", \
	RUN+="/sbin/insmod /home/gohkim/study/little_challenge/task05/drv.ko"
ACTION=="remove", ATTRS{idVendor}=="04ca", ATTRS{idProduct}=="0061", \
	RUN+="/sbin/rmmod drv"

LABEL="task05_end"
